# deploy/bookprepper-healthcheck.service
[Unit]
Description=BookPrepper health check

[Service]
Type=oneshot
ExecStart=/opt/bookprepper/deploy/bookprepper-healthcheck.sh
